# dreamia.service
[Unit]
Description=Dreamia server running LLMs for NPCs.

[Service]
Type=simple
User=USER
Group=USER
ExecStart=/bin/bash /PATH/TO/DREAMIA/start_linux.sh
WorkingDirectory=/PATH/TO/DREAMIA/

[Install]
WantedBy=multi-user.target